# services/pyload/run
#!/bin/bash
umask 0002

#Checking if config files exists (if not, copy default ones)
if  [ ! -f "/config/files.db" ] || [ ! -f "/config/files.version" ] || [ ! -f "/config/plugin.conf" ] || [ ! -f "/config/pyload.conf" ]
then
        mkdir -p /config
        chmod 777 /config

        mv /tmp/pyload-config/* /config/
fi

if [ -f "/config/pyload.pid" ]
then
        rm -f /config/pyload.pid
fi

#Changing rights if needed
chown -v abc:abc -R /opt/pyload
chown -v abc:abc -R /config

#Starting pyload
exec /sbin/setuser abc /usr/bin/python /opt/pyload/pyLoadCore.py
